# Min.conf
[program:Min]
command=python3 /home/curator4/Min/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/Min/test.err.log
stdout_logfile=/home/curator4/Min/test.out.log